# 100-dp-switch.rules
ACTION=="change" KERNEL=="extcon0" RUN+="/usr/bin/dock-script.sh"
ACTION=="add" KERNEL=="extcon0" RUN+="/usr/bin/dock-script.sh"
